# http-server-manager.service
[Unit]
Description=HTTP Server Manager
Documentation=https://github.com/tallen5431/http_server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/user/HTTP_Server
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="NODE_ENV=production"
Environment="CONFIG_FILE=/home/user/HTTP_Server/config.json"
Environment="PORT=3000"
# PROJECTS_DIR is intentionally unset here so it falls back to server.js's
# default of <WorkingDirectory>/projects. Only set it explicitly if projects
# should live outside this checkout.
ExecStart=/usr/bin/node /home/user/HTTP_Server/server.js
Restart=on-failure
RestartSec=10s

# Graceful shutdown
KillMode=process
KillSignal=SIGINT
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
